# grafana.service
[Unit]
Description=service to start the Grafana monitoring and alerting software
After=network.target

[Service]
User=ubuntu
Type=simple
Restart=always
RestartSec=1s
WorkingDirectory=/home/ubuntu/
ExecStart=/home/ubuntu/grafana-7.0.1/bin/grafana-server

[Install]
WantedBy=default.target

